# systemd/user/seven-context-observer.service
[Unit]
Description=SevenOS semantic context observer
PartOf=sevenos-session.target
After=seven-daemon.service graphical-session.target

[Service]
Type=simple
Environment=SEVENOS_ROOT=%h/Code/OS/SevenOS
Environment=SEVENOS_CONTEXT_INTERVAL=60
PassEnvironment=WAYLAND_DISPLAY DISPLAY XDG_CURRENT_DESKTOP XDG_SESSION_DESKTOP XDG_SESSION_TYPE
ExecStart=%h/.local/bin/seven-daemon observe-loop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=sevenos-session.target
