# systemd/user/sevenos-shell-experience.service
[Unit]
Description=SevenOS shell experience warmup
PartOf=sevenos-session.target
After=sevenos-wallpaper.service sevenos-waybar.service

[Service]
Type=oneshot
Environment=SEVENOS_EXPERIENCE_SILENT_WARMUP=1
ExecStart=%h/.local/bin/seven experience warmup

[Install]
WantedBy=sevenos-session.target
