# systemd/user/sevenos-notifications.service
[Unit]
Description=SevenOS notification surface
PartOf=sevenos-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/env bash -lc 'command -v mako >/dev/null 2>&1 || exit 0; exec mako'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=sevenos-session.target
