# systemd/user/sevenos-session.target
[Unit]
Description=SevenOS desktop session services
Documentation=https://github.com/seven-v73/SevenOS

[Install]
WantedBy=default.target
